# CloverPackage/sym/package/Clover/rc.scripts.core.pkg/Scripts/Tools/rc.files
etc/rc.boot.d/10.save_and_rotate_boot_log.local
etc/rc.boot.d/20.mount_ESP.local
etc/rc.boot.d/70.disable_sleep_proxy_client.local.disabled
etc/rc.clover.lib
etc/rc.shutdown.d/80.save_nvram_plist.local
Library/Application Support/Clover/CloverDaemon
Library/Application Support/Clover/CloverDaemon-stopservice
Library/LaunchDaemons/com.projectosx.clover.daemon.plist
